# deb/src/deb/init/auth.conf
description "TheSett Auth Service"
author "Rupert Smith"

start on runlevel [2345]
stop on runlevel [!2345]

setuid auth
chdir /opt/auth

script
    ./run.sh
end script
